# containers/rootfs/etc/s6-overlay/s6-rc.d/hubproxy/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
s6-setuidgid www-data
cd /opt/hubproxy/
/opt/hubproxy/hubproxy